# deploy/n8n/systemd/coa-n8n-backup.service
[Unit]
Description=COA NEWS n8n encrypted backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/n8n/backup.sh
User=root
